# linux/debian/configs/update.service
[Unit]
Description=Update system packages

[Service]
Type=oneshot
ExecStart=%h/Workspace/os_bootstrap/linux/debian/update.sh
